# Stable Version 1/static/assets/js/etc/init/print-image-id.conf
exec /bin/cat /etc/image-id
console output
start on startup
