# Stable Version 1/static/assets/js/etc/init/update-motd.conf
start on stopped rc RUNLEVEL=[2345]
exec /usr/sbin/update-motd
