# Stable Version 1/static/assets/js/etc/NetworkManager/dispatcher.d/10-sendmail
#!/bin/sh

if [ "$2" = "up" ] || [ "$2" = "down" ]; then
	/sbin/service sendmail reload || :
fi
